# systemd/example.service.in
[Unit]
Description=Example Service

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/x_PROJECT_NAME_x
Restart=on-failure

[Install]
WantedBy=multi-user.target
